# mqttx-web/rootfs/etc/s6-overlay/s6-rc.d/mqttx/run
#!/usr/bin/env bashio
# ==============================================================================
# /etc/cont-init.d/banner.sh

# Start NGiNX
#/etc/services.d/nginx/run &

## Start MQTTX Web
cd /app/web/
/app/web/http-server -p 8080 &
